# debian/cleverdb-agent.upstart
description "CleverDB Agent"

respawn
respawn limit unlimited

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

exec /usr/bin/cleverdb-agent
